# systemd/cliphist.service
[Unit]
Description=Clipboard History Daemon (GNOME Wayland)
Documentation=https://github.com/you/cliphist-x11
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/clipdaemon.py
Restart=on-failure
RestartSec=3s
Environment=GDK_BACKEND=wayland
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=GNOME
PassEnvironment=DISPLAY WAYLAND_DISPLAY XDG_RUNTIME_DIR

[Install]
WantedBy=graphical-session.target
